# dist/power-profiles-switcher.service
[Unit]
Description=Power Profiles Switcher
Requires=dbus.socket
After=power-profiles-daemon.service
After=upower.service

[Service]
Type=simple
ExecStart=%h/bin/power-profiles-switcher

[Install]
WantedBy=default.target
